# .config/systemd/user/morse-notifier.service
[Unit]
Description=Morse Notifier
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

[Service]
Type=exec
ExecStart=/home/bczhc/.config/systemd/user/morse-notifier
Restart=on-failure

[Install]
WantedBy=graphical-session.target
